# drbd/scripts/drbd
#!/bin/sh
#
# chkconfig: 345 30 70
# description: Loads and unloads the drbd module
#	       

# Source function library.
. /etc/rc.d/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

# Check for drbdsetup
[ -x /usr/sbin/drbdsetup ] || exit 1

cd /etc/sysconfig/drbd
DEVICES=`ls drbd[0-9] 2>/dev/null`
HOST=`hostname -s`

# See how we were called.
case "$1" in
  start)
        action "Loading DRBD module" /sbin/modprobe drbd
	for D in $DEVICES; do
	    . $D
	    if [ ${SIZE} != "" ]; then
		SIZE_OP="-d ${SIZE}";
	    fi
	    if [ ${HOST} = `echo ${MASTER_NODE} | cut -d: -f1` ]; then
		action "Configuring DRBD resource %s" $D \
		    /usr/sbin/drbdsetup $MASTER_DEVICE $MASTER_PARTITION \
		    $PROTOCOL $MASTER_NODE $SLAVE_NODE $SIZE_OP $OPTIONS
		action "Waiting for DRBD resource %s to resynchronize" $D \
		    /usr/sbin/drbdsetup $MASTER_DEVICE WAIT
	    elif [ ${HOST} = `echo ${SLAVE_NODE} | cut -d: -f1` ]; then
		action "Configuring DRBD resource %s" $D \
		    /usr/sbin/drbdsetup $SLAVE_DEVICE $SLAVE_PARTITION \
		    $PROTOCOL $SLAVE_NODE $MASTER_NODE $SIZE_OP $OPTIONS
		action "Waiting for DRBD resource %s to resynchonize" $D \
		    /usr/sbin/drbdsetup $SLAVE_DEVICE WAIT
	    else
		echo This machine is neither master nor slave!
		exit 1
	    fi
	done
	;;
  stop)
	action "Unloading DRBD module" /sbin/modprobe -r drbd
	;;
  restart)
	$0 stop
	$0 start
	;;
  status)
	if [ -r /proc/drbd ]; then
		cat /proc/drbd
	else
		echo "DRBD module not loaded."
	fi
	;;
  *)
	echo "Usage: $0 {start|stop|restart|status}"
	exit 1
esac

exit 0
