# debian/drbd0.7-utils.postrm
#!/bin/sh

if [ "$1" = "purge" ] ; then
	update-rc.d drbd remove >/dev/null
fi
#DEBHELPER#
