# debian/drbd8-utils.postrm
#!/bin/sh

case $1 in purge|remove|deconfigure)
	update-rc.d drbd remove >/dev/null
	# remove compat symlinks installed by .postinst
	rm -f /sbin/drbdsetup /sbin/drbdmeta /sbin/drbdadm
	;;
esac

#DEBHELPER#
