# debian/drbd8-utils.prerm
#!/bin/sh

# This script silently ignores the initscript return code if we remove
# or deconfigure the package or carefully pay attention to it if we
# upgrade the package.

# Cyril Bouthors <cyril@bouthors.org>
# Wed Feb 16 21:01:11 CET 2005

set -e

if [ -x "/etc/init.d/drbd" ]
then
    case "$1" in
	remove)
	    rm -f /etc/udev/rules.d/65_drbd.rules
	    # fall through
	deconfigure)
	    /usr/sbin/invoke-rc.d drbd stop || true
	    ;;

	upgrade|failed-upgrade)
	    /usr/sbin/invoke-rc.d drbd stop
	    ;;

	*)
	    echo "prerm called with unknown argument \`$1'" >&2
	    exit 0
	    ;;
    esac
fi

#DEBHELPER#
